# supervisord.conf
[supervisord]
nodaemon = true

[program:apache]
command=apachectl -D "FOREGROUND" -k start
redirect_stderr=true

[include]
files = supervisor/conf.d/*.conf
